# brightline/tools/com.brightline.r2watcher.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.brightline.r2watcher</string>
	<key>ProgramArguments</key>
	<array>
		<string>/bin/bash</string>
		<string>/Users/kiril/brightline/tools/start-r2watcher.sh</string>
	</array>
	<key>WorkingDirectory</key>
	<string>/Users/kiril/brightline/tools</string>
	<key>RunAtLoad</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/Users/kiril/Library/Logs/brightline-r2watcher.log</string>
	<key>StandardErrorPath</key>
	<string>/Users/kiril/Library/Logs/brightline-r2watcher.error.log</string>
</dict>
</plist>
